# home/users/ejvend/config/autostart/Qlipper.desktop
[Desktop Entry]
Exec=qlipper
Name=Qlipper
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
